# docker/cron-docker/periodic/hourly/restart_worker
#!/bin/sh
set -e
# find worker and restart the container(s)
docker ps --filter label=com.docker.compose.service=worker --quiet | xargs --no-run-if-empty docker restart
