# docker/cron-docker/periodic/weekly/docker_image_prune
#!/bin/sh
set -e
# remove unused images, --all to remove images not used by containers, --force to skip confirm prompt
docker image prune --all --force
